# workspace/crontab/laradock
* * * * * laradock /usr/bin/php /var/www/artisan schedule:run >> /dev/null 2>&1
*/10 * * * * laradock /usr/bin/php /var/www/artisan make:controller test
